# thonnycontrib/ev3dev/res/rpycd.service
[Unit]
Description=Classic rpyc server 


[Service]
Type=simple
ExecStart=/usr/bin/rpyc_classic.py

[Install]
WantedBy=multi-user.target

